# contrib/cam-reverse.service
# Example systemd service for the cam-reverse HTTP server.
#
# Install:
#   sudo cp contrib/cam-reverse.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cam-reverse
#
# Adjust User, WorkingDirectory, ExecStart and the config path to your setup.

[Unit]
Description=cam-reverse IP camera HTTP gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cam-reverse
Group=cam-reverse
WorkingDirectory=/opt/cam-reverse
# Path to the cam-reverse.py entry point inside your virtualenv
ExecStart=/opt/cam-reverse/.venv/bin/cam-reverse.py http_server --config_file /etc/cam-reverse/config.yml
Restart=on-failure
RestartSec=5

# Hardening (loosen if it interferes with your deployment)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Allow writing the config file back from the web UI; drop if you don't use that
ReadWritePaths=/etc/cam-reverse

[Install]
WantedBy=multi-user.target
